# webcrawler/systemd/napkin_collector.timer
[Unit]
Description=Run the napkind collector webcrawler every 15 min.

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
