# webcrawler/systemd/napkin_collector.service
[Unit]
Description=Launch the napkin collector webcrawler

[Service]
Type=oneshot
ExecStart=/usr/bin/php56 /var/www/webcrawler/driver.php
